# files/ubuntu/sv/update-notifier/run
#!/bin/sh

service=update-notifier

. ../common.sh

depends_on awesome

exec /usr/bin/update-notifier
