# files/ubuntu/sv/print-applet/run
#!/bin/sh

exec /usr/bin/system-config-printer-applet
