# files/sv/awesome/run
#!/bin/sh

service=awesome
enable_logging=true

. ../common.sh

export wpsetters=feh
exec awesome
