# files/ubuntu/profile.d/chruby.sh
export PATH="/usr/local/chruby/bin:$PATH"
export MANPATH="/usr/local/chruby/share/man:$MANPATH"
