# files/ubuntu/profile.d/ruby-install.sh
export PATH="/usr/local/ruby-install/bin:$PATH"
export MANPATH="/usr/local/ruby-install/share/man:$MANPATH"
