# files/ubuntu/sv/kupfer/run
#!/bin/sh

service=kupfer
enable_logging=true

. ../common.sh

exec /usr/bin/kupfer --no-splash
