# files/sv/gnome-keyring/run
#!/bin/sh

exec /usr/bin/gnome-keyring-daemon --start --foreground
